# roles/img-panda/files/supervisord/imgpanda_app.conf
[program:img_panda_app]
command=python /tmp/imgpanda_app/app.py
autostart=true
autorestart=true
stderr_logfile=/var/log/imgpanda_app.err.log
stdout_logfile=/var/log/imgpanda_app.out.log
