# roles/smart-panda/files/supervisord/smart_panda.conf
[program:smart_panda]
command=python /tmp/smart_panda/app.py
autostart=true
autorestart=true
stderr_logfile=/var/log/smart_panda.err.log
stdout_logfile=/var/log/smart_panda.out.log
